# tools/ops-alerts/systemd/ops-alerts.timer
[Unit]
Description=Home SOC - Run operational alerts daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
