# tools/ops-alerts/systemd/ops-alerts.service
[Unit]
Description=Home SOC - Operational alerts runner (P1)
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=socadmin
WorkingDirectory=/home/socadmin/soc-cases
# Load Indexer credentials (file perms 600) and run checks
ExecStart=/bin/bash -lc 'set -a && source /home/socadmin/.secrets/mini-soc.env && set +a && /home/socadmin/soc-cases/tools/ops-alerts/run_ops_alerts.sh'
